# production/quantum_notice_app/deploy/systemd/start_for_weekday.service
[Unit]
Description=Start weekday services for the week

[Service]
Type=oneshot
ExecStart=/home/norma/quantum_notice_app/scripts/manage_services.sh start
