# production/quantum_notice_app/deploy/systemd/stop_for_weekend.service
[Unit]
Description=Stop weekday services for the weekend

[Service]
Type=oneshot
ExecStart=/home/norma/quantum_notice_app/scripts/manage_services.sh stop
